# overlay/elvees/common/usr/lib/systemd/system/alsa-restore.service.d/mfbsp-i2s.conf
# this will start asound-restore after mfbsp-i2s sound card appears
[Unit]
Requires=sys-devices-platform-aux_sound-sound-card1-controlC1.device
After=sys-devices-platform-aux_sound-sound-card1-controlC1.device
